# themes/root/etc/uci-defaults/luci-theme-icstheme
	#!/bin/sh
	uci batch <<-EOF
	        set luci.themes.MyTheme=/luci-static/icstheme
	        commit luci
	EOF
exit 0